# systemd/user/app-lxqt\x2dpowermanagement@autostart.service.d/override.conf
[Unit]
ConditionEnvironment=!LXQT_SESSION_MANAGER=uwsm
